# network_information_source/run/plan_out.cron
PYTHONPATH=$PATHONPATH:/u01/isi/CenterProject/network_information_source

# 发现/验证任务执行
*/10 * * * * cd /u01/isi/CenterProject/network_information_source/run && bash main.sh start deploy_find.py
*/10 * * * * cd /u01/isi/CenterProject/network_information_source/run && bash main.sh start deploy_verify.py

* */12 * * * cd /u01/isi/CenterProject/network_information_source/run && bash clear_log.sh